# notifiarr/rootfs/etc/s6-overlay/s6-rc.d/notifiarr/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant  Add-on: Notifiarr
# Runs Notifiarr
# ==============================================================================

bashio::log.info 'Starting Notifiarr...'

exec /opt/bin/Notifiarr -nobrowser -data=/config
